# exam_1/ansible/roles/worker/files/datanode.service
[Unit]
Description=Hadoop datanode
After=network.target network-online.target

[Service]
User=hdfs
Group=hadoop
Type=fork
ExecStart=/usr/local/hadoop/current/bin/hdfs datanode
Restart=on-failure

[Install]
WantedBy=multi-user.target
